# packages/core/src/main/product/app/system/cron/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
*/5 * * * * @{product.bin}/cron.5mins.sh
