# packages/core/src/main/product/bin/profile.d/core.sh
#!/bin/sh
#
# Core profile configuration.
# By Fabien CRESPEL <fabien@crespel.net>
#

umask 022
